# files/etc/hotplug.d/iface/50-ratelimit
/etc/init.d/ratelimit running || exit 0
[ "$ACTION" = "ifup" ] || exit 0
ubus call ratelimit reload
